# tools/acm_merge.service
[Unit]
Description=Accomplishments Exporter merge

[Service]
Type=oneshot
# CONFIGURE: install dir (/opt/acm-exporter), lua path (/usr/bin/lua), and cluster
# root + output path (/path/to/cluster) must be edited for your deployment.
Environment=LUA_PATH=/opt/acm-exporter/scripts/?.lua;/opt/acm-exporter/tools/?.lua;;
ExecStart=/usr/bin/lua /opt/acm-exporter/tools/acm_merge.lua --root /path/to/cluster --out /path/to/cluster/acm_export.json
